# service/.local/service/ssh-agent/run
#!/bin/sh
exec ssh-agent -D > "$HOME/.ssh/ssh-agent-env"
